# slowproxy.conf
[program:slowproxy]
directory=/home/pradeep/scripts
command=ruby /home/pradeep/scripts/slowproxy.rb -f /home/pradeep/scripts/slowsites.txt -p 5310 -d 30
user=pradeep
